# deploy/cron/reeflex-nilm-daily.cron
30 2 * * * cd /opt/reeflex-nilm && /usr/bin/docker compose run --rm -T reeflex-nilm-daily >> /var/log/reeflex-nilm-daily.log 2>&1
